# totoro-device/scripts/S99totoro-device
#!/bin/sh

### BEGIN INIT INFO
# Provides:          totoro-device
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs $network
# Default-Start:     S
# Default-Stop:      0 6
# Short-Description: Totoro Device (NWCT) service
### END INIT INFO

DAEMON=/root/totoro-device
NAME=totoro-device
PIDFILE=/var/run/${NAME}.pid
LOGFILE=/var/log/${NAME}.log

# 持久化配置（重启不丢）
export NWCT_CONFIG_PATH=/etc/nwct/config.json
export NWCT_DB_PATH=/var/nwct/devices.db
export NWCT_LOG_DIR=/var/log/nwct

start() {
  echo "Starting ${NAME}..."
  mkdir -p /var/run /var/log/nwct /var/nwct

  if [ ! -x "${DAEMON}" ]; then
    echo "Missing executable: ${DAEMON}"
    return 1
  fi

  # 已在运行则跳过
  if [ -f "${PIDFILE}" ]; then
    pid=$(cat "${PIDFILE}" 2>/dev/null)
    if [ -n "${pid}" ] && kill -0 "${pid}" 2>/dev/null; then
      echo "${NAME} already running (pid=${pid})"
      return 0
    fi
  fi

  # 后台启动（display=true 驱动 fb0）
  nohup "${DAEMON}" -display=true >>"${LOGFILE}" 2>&1 &
  echo $! >"${PIDFILE}"
  sleep 1

  pid=$(cat "${PIDFILE}" 2>/dev/null)
  if [ -n "${pid}" ] && kill -0 "${pid}" 2>/dev/null; then
    echo "${NAME} started (pid=${pid})"
    return 0
  fi

  echo "${NAME} failed to start, see ${LOGFILE}"
  return 1
}

stop() {
  echo "Stopping ${NAME}..."
  if [ ! -f "${PIDFILE}" ]; then
    echo "${NAME} not running (no pidfile)"
    return 0
  fi

  pid=$(cat "${PIDFILE}" 2>/dev/null)
  if [ -z "${pid}" ]; then
    rm -f "${PIDFILE}"
    return 0
  fi

  kill "${pid}" 2>/dev/null || true
  # 等待退出
  for i in 1 2 3 4 5 6 7 8 9 10; do
    if kill -0 "${pid}" 2>/dev/null; then
      sleep 1
    else
      break
    fi
  done
  if kill -0 "${pid}" 2>/dev/null; then
    kill -9 "${pid}" 2>/dev/null || true
  fi

  rm -f "${PIDFILE}"
  echo "${NAME} stopped"
}

status() {
  if [ -f "${PIDFILE}" ]; then
    pid=$(cat "${PIDFILE}" 2>/dev/null)
    if [ -n "${pid}" ] && kill -0 "${pid}" 2>/dev/null; then
      echo "${NAME} running (pid=${pid})"
      return 0
    fi
  fi
  echo "${NAME} not running"
  return 3
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    start
    ;;
  status)
    status
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac


